# 99-userdev-input.rules
# set rtc so that user dextra can access without root

KERNEL=="event*", SUBSYSTEM=="input", RUN+="/usr/bin/setfacl -m u:dextra:rw /dev/rtc"
